# deploy/ruby-updated.service
[Unit]
Description=Privileged Ruby updater (queue consumer)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ruby-updated --consume
TimeoutStartSec=15min
