# deploy/rubyvision.service
[Unit]
Description=Ruby camera -> Hailo inference -> annotated frame service
After=network.target
# Never enter the systemd start-rate-limiter: this service may legitimately
# crash-loop while a camera/HEF is being set up, and Restart=always must keep
# retrying. StartLimitIntervalSec/StartLimitBurst are [Unit]-section directives
# (ignored if placed under [Service]).
StartLimitIntervalSec=0

[Service]
User=michael
Group=michael
# video: V4L2 / camera nodes; render: Hailo /dev/hailo* DRI access.
SupplementaryGroups=video render
# rubyvision is imported via PYTHONPATH rather than a pip editable install: the
# PEP-660 editable finder bakes in an absolute path, and the OTA updater stages
# each release in a throwaway worktree -- so an editable install can end up
# pointing at a pruned release dir (dangling -> "No module named rubyvision").
# /home/michael/hud is the symlink to the active release's hud/; vision/ is its
# sibling, so /home/michael/hud/../vision always resolves to the LIVE release's
# vision package and auto-tracks OTA updates (restart rubyvision after an apply).
Environment=PYTHONUNBUFFERED=1 PYTHONPATH=/home/michael/hud/../vision
ExecStart=/home/michael/ruby-env/bin/python -m rubyvision --models /home/michael/vision/models --shm /dev/shm/rubyvision --fps 15
Restart=always
RestartSec=2
Nice=5

[Install]
WantedBy=multi-user.target
